# installer/linux/codeshelf.conf
description "start and stop the codeshelf server"

start on runlevel [2345]
stop on runlevel [!2345]

env ROOT_DIR=/opt/codeshelf/engine
env JAVA_HOME=/usr/lib/jvm/java-7-oracle
env JAVA_OPTS="-Xms1536M -Xmx1536M"
env USER=codeshelf
env LOGFILE=/home/codeshelf/.CodeShelf/logs/daemon.log

#setuid codeshelf

#expect fork
respawn
respawn limit 10 5
chdir /opt/codeshelf/engine

exec su -s /bin/sh -c 'exec "$0" "$@"' codeshelf -- /usr/bin/java -Dconfig.properties=/etc/codeshelf/server.config.properties -jar codeshelf.jar > $LOGFILE 2>&1

#exec /usr/bin/java -Xms1536M -Xmx1536M -jar minecraft_server.jar nogui > $LOGFILE 2>&1

#exec start-stop-daemon --start --exec $JAVA_HOME/bin/java -- $JAVA_OPTS -jar $ROOT_DIR/minecraft_server.jar nogui > $LOGFILE 2>&1

